# etc/services.d/main_run/run
#!/command/with-contenv bash
while true; do sleep 60; done;
